# dotfiles/Documents/Services/change_wallpaper.timer
[Unit]
Description=Periodically change wallpaper

[Timer]
OnStartupSec=120s
OnUnitActiveSec=120s

[Install]
WantedBy=timers.target

