# dotfiles/Documents/Services/clean_downloads.timer
[Unit]
Description=Periodically move old downloads to tmp folder

[Timer]
OnStartupSec=5min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target

